# root/etc/services.d/darkhttpd/run
#!/usr/bin/with-contenv bash

#export HOME="~app"
#exec s6-setuidgid root "/config/start_darkhttpd.sh"
exec /config/start_darkhttpd.sh
